# contrib/upstart.conf
description "simpledoge"
 
start on (filesystem) 
stop on runlevel [016] 
 
respawn 
console log 
setuid dogecoin 
setgid dogecoin 
chdir /home/dogecoin/simpledoge 
 
exec /home/dogecoin/venv/bin/gunicorn simpledoge.uwsgi_entry:app 
